# utility/cassandra.conf
[program:start_cassandra]
command=/bin/run.sh
